# 3.1_data_sampling/src/set_buffer_size.service
[Unit]
Description=Set buffer size

[Service]
ExecStart=/boot/set_buffer_size.sh

[Install]
WantedBy=multi-user.target
